# mcp_reporter/rootfs/etc/services.d/mcp-report/run
#!/usr/bin/with-contenv sh
python -m mcp_reporter_server
